# support/com.sovereign.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sovereign.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>LAUNCH_SCRIPT_PLACEHOLDER</string>
    </array>
    <key>WorkingDirectory</key>
    <string>REPO_DIR_PLACEHOLDER</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>STDOUT_LOG_PLACEHOLDER</string>
    <key>StandardErrorPath</key>
    <string>STDERR_LOG_PLACEHOLDER</string>
</dict>
</plist>
